# luge/pi/sledc.service
[Unit]
Description=Sled Client
Documentation="TBD"

[Service]
ExecStart=/root/sledc run
Type=simple
Restart=on-failure
RestartSec=5
